# headphones-plug-detector.service
[Unit]
Description=Headphones plug detector
After=dbus.socket

[Service]
ExecStart=@EXECUTABLE_PATH@

[Install]
WantedBy=default.target
